# discord-bot.service
[Unit]
Description=Discord Spam Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/discord_spam_bot
Environment=PATH=/home/ubuntu/discord_spam_bot/venv/bin
ExecStart=/home/ubuntu/discord_spam_bot/venv/bin/python bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 